# templates/default/init.sh.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides: niet
# Defalt-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: Niet global service configuration manager
### END INIT INFO

# This config file is a template managed via chef - don't edit it

start() {
  <% @processes.each do |user, locations| %>
      <% locations.each do |loc| %>
        <% if loc.include?("\*") %>
    for f in <%= loc %>
    do
      su --login <%= user %> $f
    done
        <% else %>
    su --login <%= user %> <%= loc %>
        <% end %>
      <% end %>
    <% end %>
}
stop() {
<% @processes.keys.each do |user| %>
    while pkill -QUIT -u <%= user %> niet; do sleep 1; done
<% end %>
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart)
        # Graceful restart
        stop && start
        ;;
  kill)
        while pkill -KILL niet; do sleep 1; done
        ;;
  *)
    echo "Usage: $0 {start|stop|restart|kill}" >&2
        exit 1
        ;;
  esac

  exit 0
